# deploy/systemd/visor-refresh-margenes.timer
[Unit]
Description=Visor - timer diario de refresh de la vista materializada de margenes
Documentation=man:systemd.timer(5)

[Timer]
# Todos los dias a las 06:00 hora del sistema. Asume que la ETL nocturna
# termina antes; ajustar si la ETL corre mas tarde para evitar refrescar
# antes de que llegue la data nueva.
OnCalendar=*-*-* 06:00:00
# Si la VM estaba apagada en ese momento, ejecutar al arrancar
Persistent=true
# Dispersa la ejecucion para evitar picos
RandomizedDelaySec=300
Unit=visor-refresh-margenes.service

[Install]
WantedBy=timers.target
